# com.insightbridge.ibkr.paper.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.insightbridge.ibkr.paper.gateway</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/tongyin/Desktop/InsightBridge_Financial_Models_Latest/start_ibkr_paper_gateway.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/tongyin/Desktop/InsightBridge_Financial_Models_Latest/ibkr_runtime/logs/launchd_stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/tongyin/Desktop/InsightBridge_Financial_Models_Latest/ibkr_runtime/logs/launchd_stderr.log</string>
</dict>
</plist>
