# com.insightbridge.dukascopy.fx.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.insightbridge.dukascopy.fx.bridge</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/anaconda3/bin/python3</string>
    <string>-u</string>
    <string>/Users/tongyin/Desktop/InsightBridge_Financial_Models_Latest/03_FX_AUD_NZD_EUR_GBP/backend/server.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/tongyin/Desktop/InsightBridge_Financial_Models_Latest/03_FX_AUD_NZD_EUR_GBP/backend</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>DUKASCOPY_FX_PAIRS</key>
    <string>AUD/USD,NZD/USD,EUR/USD,USD/JPY,GBP/USD,AUD/JPY,NZD/JPY</string>
    <key>PYTHONPATH</key>
    <string>/Users/tongyin/Desktop/InsightBridge_Financial_Models_Latest</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/tongyin/Desktop/InsightBridge_Financial_Models_Latest/reports/dukascopy_bridge/launchd_stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/tongyin/Desktop/InsightBridge_Financial_Models_Latest/reports/dukascopy_bridge/launchd_stderr.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
